# scripts/krossboard-data-processor-consolidator.timer
[Unit]
Description=Run Krossboard clusters data collector periodically

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target