# scripts/krossboard-data-processor-consolidator.timer
[Unit]
Description=Krossboard Analytics Consolidator Trigger

[Timer]
OnCalendar=*:0/5

[Install]
WantedBy=timers.target