# scripts/krossboard-data-processor-collector.timer
[Unit]
Description=Krossboard Data Collector Trigger

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target